# etc/systemd/system/aiq-backend.service
[Unit]
Description=AIQ Backend Docker Container Service

Ensure Docker is running before attempting to start this service
Requires=docker.service
After=docker.service

[Service]

Always restart the application if it fails or crashes
Restart=always

Run the service as the root user for Docker access
User=root

--- !!! IMPORTANT: Replace 'your_username' with your actual image tag !!! ---
ExecStart=/usr/bin/docker run --name aiq_backend -p 5000:5000 your_username/aiq-backend:latest

Command to gracefully stop the container
ExecStop=/usr/bin/docker stop aiq_backend

Command to remove the container instance after stopping (allows clean restarts)
ExecStopPost=/usr/bin/docker rm aiq_backend

[Install]

Enable the service to start automatically when the server boots up
WantedBy=multi-user.target