# etc/systemd/system/aiq-backend.service
[Unit]
Description=AIQ Backend Docker Container Service
Requires=docker.service
After=docker.service

[Service]

# Set the restart policy to always restart the container if it fails or stops

Restart=always
User=root

# CRITICAL FIX: The ExecStart command is hardcoded to force Gunicorn to use

# a single worker (--workers 1). This permanently resolves the Socket.IO

# connection issue that the frontend was experiencing.

ExecStart=/usr/bin/docker run --name aiq_backend -p 5000:5000 aiq-backend-prod:latest gunicorn --workers 1 --worker-class eventlet -b 0.0.0.0:5000 'aiq-backend:app'

# Commands to gracefully stop and clean up the container

ExecStop=/usr/bin/docker stop aiq_backend
ExecStopPost=/usr/bin/docker rm aiq_backend

[Install]
WantedBy=multi-user.target